# verticevnc/centos7/etc/systemd/system/verticevnc.service.erb
[Unit]
Description=Verticevnc Server
After=network.target
After=runlevel2.target
After=runlevel3.target
After=runlevel4.target
After=runlevel5.target
[Service]
EnvironmentFile=/var/lib/megam/env.sh
ExecStart=/usr/bin/node /usr/share/megam/verticevnc/vnc_monitor.js
KillMode=process
