# verticecadvisor/xenial/etc/systemd/system/cadvisor.service.erb
[Unit]
Description=cAdvisor
After=network.target
After=runlevel2.target
After=runlevel3.target
After=runlevel4.target
After=runlevel5.target
[Service]
ExecStart=<%= @basic[:home] %>/cadvisor/bin/cadvisor -port 9999
KillMode=process
