# verticensqd/centos7/etc/systemd/system/nsqd.service
[Unit]
Description=Vertice nsqd
After=network.target
After=runlevel2.target
After=runlevel3.target
After=runlevel4.target
After=runlevel5.target
[Service]
ExecStart=/usr/share/megam/nsqd/bin/nsqd --lookupd-tcp-address=localhost:4160 -data-path=/var/lib/megam/vertice/nsqd
KillMode=process
