# verticevnc/xenial/etc/systemd/system/verticevnc.service.erb
[Unit]
Description=verticevnc
After=network.target
After=runlevel2.target
After=runlevel3.target
After=runlevel4.target
After=runlevel5.target
[Service]
WorkingDirectory=<%= @basic[:home]+"/"+@package[:package] %>
ExecStart=/usr/bin/node vnc_monitor.js
